# qubes-flatpak-system-monitor.service
[Unit]
Description=Flatpak install/uninstall monitor
[Service]
Type=simple
ExecStart=/etc/qubes-rpc/qubes.PostInstall
