# qubes-flatpak-system-monitor.path
[Unit] 
Description=Flatpak install/uninstall system monitor
ConditionPathExists=/usr/bin/flatpak
ConditionPathExists=/var/run/qubes/this-is-templatevm
[Path]
PathChanged=/var/lib/flatpak/.changed
[Install]
WantedBy=multi-user.target

